# etc/sv/sms/run
#!/bin/bash
exec 2>&1

if [ ! -e /dev/ttyUSB2 ]; then
	sleep 2
	exit 1
fi

exec gammu-smsd --config /etc/gammu-smsdrc --pid /var/run/gammu-smsd.pid
