# etc/sv/kotel/run
#!/bin/bash

exec 2>&1

echo
echo "kotel start..."

exec python -u /opt/kotel/kotel.py -c /opt/kotel/kotel.cfg
